# nodered.service
[Unit]
Description=Node-RED
After=network.target

[Service]
Type=simple
User=ubuntu
ExecStart=/usr/local/bin/node-red
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
